# debian/etc/cron.d/pbspro-execution-meta
MAILTO=vchlumsky@cesnet.cz
26-59/30 * * * * root /var/spool/pbs/mom_scripts/touch_scratch
